# home/.config/systemd/user/app-org.fcitx.Fcitx5@autostart.service.d/override.conf
[Service]
# cskk-git installs libcskk under /usr/lib/cskk; make it visible to fcitx5-cskk.
Environment=LD_LIBRARY_PATH=/usr/lib/cskk
